# scripts/systemd/disk-space-guard.service
# NOTE: replace /path/to/marveen (install dir) and /home/USER (home) with your values.
[Unit]
Description=Marveen disk-space guard (reap scratch + alert before the root fs fills)
After=marveen-dashboard.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/path/to/marveen/scripts/disk-space-guard.sh
Environment=PATH=/home/USER/.local/bin:/home/USER/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/USER
Environment=TZ=Europe/Budapest
# journal, NOT append:file -- under disk-full (this guard's whole reason to exist)
# or a missing store/, opening an append fd would fail and the guard never starts.
StandardOutput=journal
StandardError=journal
